# files/homeassistant.service
[Unit]
Description=Homeassistant

[Service]
Type=simple
ExecStart=/usr/local/bin/hass

[Install]
WantedBy=multi-user.target
